# upstart/hue.conf
# hue service

description "hue service"

# Stanzas
#

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Specify the process/command to start, e.g.
exec start-stop-daemon --start --chuid hue:hue --exec /usr/share/hue/build/env/bin/supervisor
